# out/v5-hdmi/compile_dir/target/rootfs/usr/lib/opkg/info/middleware-libs.list
/usr/lib/libmedia_utils.so
/usr/lib/libkfc_server_enc.so
/usr/bin/kfc_server_pack
/usr/lib/libcdx_stream.so
/usr/lib/lib_ise_sti.so
/usr/lib/libvencoder.so
/usr/lib/libcdx_base.so
/usr/lib/libawh264.so
/usr/lib/libcdx_common.so
/usr/lib/libz.so.1.2.8
/usr/lib/libmpp_ise.so
/usr/lib/libawh265.so
/usr/lib/libmpp_isp.so
/usr/lib/libMemAdapter.so
/usr/lib/libip_loader_soft.so
/usr/lib/libmpp_vi.so
/etc/init.d/kfc_server_pack
/usr/lib/libcedarx_aencoder.so
/usr/lib/libcedarxrender.so
/usr/lib/libcrypto.so.1.0.0
/usr/lib/libcdx_parser.so
/usr/lib/libvideoengine.so
/usr/lib/libmpp_component.so
/usr/lib/lib_ise_bi_soft.so
/usr/bin/iploader_server
/usr/lib/libasound.so.2
/usr/lib/libISP.so
/usr/lib/libvdecoder.so
/usr/lib/libVE.so
/etc/init.d/iploader_server
/usr/lib/libTinyServer.so
/usr/lib/libkfcapi_enc.so
/usr/lib/libasound.so
/usr/lib/libz.so.1
/usr/lib/libnormal_audio.so
/usr/lib/libcedarxstream.so
/usr/lib/libmedia_mpp.so
/usr/lib/libcdc_base.so
/usr/lib/libmpp_vo.so
/usr/lib/libawmjpegplus.so
/usr/lib/libcrypto.so
/usr/lib/lib_ise_mo.so
/usr/lib/libasound.so.2.0.0
/usr/lib/libz.so
/etc/cedarx.conf
/usr/lib/lib_ise_bi.so
/usr/lib/libadecoder.so
