# out/v5-hdmi/compile_dir/target/rootfs/usr/lib/opkg/info/tcp_server.list
/usr/bin/tcp_server
/usr/bin/resetIP.sh
/etc/init.d/autorun
/usr/bin/autosetting.sh
